# flip2/config/flip2.service
[Unit]
Description=FLIP2 Multi-Agent Coordination Daemon
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/flip2d --config /etc/flip2/config.yaml
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
